# conf/self-upgrade.timer
[Unit]
Description=Upgrade __APP__ regularly

[Timer]
OnCalendar=daily

[Install]
WantedBy=timers.target
